# config/deploy/templates/master.conf.erb
pre-start script
  mkdir -p <%= log %>
  chown -R <%= user %> <%= log %>
end script
start on (started schemer)

stop on (stopping schemer)
